# laravel/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue]
; --timeout=600: ワーカープロセスのタイムアウト. Job の $timeout (300s) より長く取らないと
;                qwen2.5:3b (CPU) で重いプロンプトを処理中にワーカーが SIGKILL されてしまい,
;                attempts++ で繰り返し max-attempts-exceeded を引いていた.
; --tries=3:    1 Job の最大試行回数. Job 側の public $tries が指定されていればそちらが優先.
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=600
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Laravel スケジューラ: bootstrap/app.php の ->withSchedule(...) を 60 秒ごとに評価。
; この中で 5 分おきに mail:fetch を起動 (POP3/IMAP からメールをローカル取得)。
; supervisord が常駐するため cron は不要。
[program:scheduler]
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
